# raspberry_pi/service/museum_service.service
[Unit]
Description=Museum Automation System - Bulletproof Edition
After=network.target network-online.target mosquitto.service
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/admin/Documents/GitHub/museum-system/raspberry_pi/main.py
ExecStart=/usr/bin/python3 /home/admin/Documents/GitHub/museum-system/raspberry_pi/watchdog.py
WorkingDirectory=/home/admin/Documents/GitHub/museum-system/raspberry_pi

# Logging with rotation
StandardOutput=inherit
StandardError=inherit

# BULLETPROOF RESTART SETTINGS
Restart=always
RestartSec=5
RestartForceExitStatus=1 2 3 4 5 6 7 8 9 130 131 132 133 134 135 136 137 138 139 255

# Watchdog settings - restart if no heartbeat for 30 seconds
WatchdogSec=30
NotifyAccess=main

# Resource limits to prevent system overload
MemoryLimit=512M
CPUQuota=50%

# User and permissions
User=admin
Group=admin
SupplementaryGroups=gpio,audio,dialout

# Environment
Environment=PYTHONPATH=/home/admin/Documents/GitHub/museum-system/raspberry_pi
Environment=PYTHONUNBUFFERED=1

# Kill settings - force kill if graceful shutdown takes too long
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
